# template/bak.crontab
# Crontab to run borgmatic once an hour
RANDOM_DELAY=20
03 * * * * /usr/bin/borgmatic -c /home/user/.ssh/backups/borgmatic.cfg
